# cron/crontab
* * * * * root /usr/local/bin/cdr-rotate >> /var/log/logger/cdr-rotate.log 2>&1
