# cron/crontab
* * * * * root /usr/local/bin/cdr-rotate
